# archive.ini
[program:archive]
command=/data/app/bimax/flow.sh archiveAll
directory=/data/app/bimax
autostart=true
autorestart=true
stdout_logfile=/var/log/archive.log
stderr_logfile=/var/log/archive.log
redirect_stderr=true

