# reduce.ini
[program:reduce]
command=/data/app/bimax/flow.sh reduceAll
directory=/data/app/bimax
autostart=true
autorestart=true
stdout_logfile=/var/log/reduce.log
stderr_logfile=/var/log/reduce.log
redirect_stderr=true

